# 99-drake.rules
KERNEL=="cpu[0-9]*", SUBSYSTEM=="cpu", RUN+="$destdir/drake_permissions $kernel"
